# misc/podwatch.service
[Unit]
Description=Podwatch updating service

[Service]
Type=oneshot
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/podwatch
