# misc/podwatch.service
[Unit]
Description=Podwatch updating service

[Service]
Type=oneshot
ExecStart=/usr/bin/podwatch
